# systemd/loopover-docker-prune@.timer.example
# Pairs with loopover-docker-prune@.service.example. See that file for install instructions -- this is
# also a templated (instantiated) unit, enabled per instance: `loopover-docker-prune@<instance>.timer`.

[Unit]
Description=Run loopover-docker-prune@%i.service daily

[Timer]
OnCalendar=daily
Persistent=true
# Spread the run over a window instead of firing at exactly midnight on every host.
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
